# station_agent/station-agent.service
[Unit]
Description=OE5XRX Station Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m station_agent
WorkingDirectory=/var/lib/station-agent
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
